# configs/systemd/.config/systemd/user/xorg.target
[Unit]
Description=Xorg environment
BindsTo=graphical.target
